# common/ubuntu/server/tailscale/host/tailscale-host-routing.service
[Unit]
Description=Tailscale subnet routing via host network mode container
After=docker-compose-auto.service
Requires=docker-compose-auto.service
StartLimitBurst=5
StartLimitIntervalSec=30

[Service]
# Service Type
# oneshot: runs once and exits
# RemainAfterExit: systemd considers service active after script completes
Type=oneshot
RemainAfterExit=yes

# Execution Commands
# Wait for Tailscale container to create tailscale0 interface
ExecStartPre=/bin/sleep 5
ExecStart=/usr/local/bin/tailscale-host-routing.sh

# Cleanup on service stop
# Remove nftables rule
ExecStop=/usr/sbin/nft delete rule ip raw PREROUTING iifname "tailscale0" accept 2>/dev/null || true

# Restart Behavior
# Restart only on failure, not on clean stop
Restart=on-failure
RestartSec=5s

# Logging
# Send all output to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=tailscale-host-routing

[Install]
WantedBy=multi-user.target
